# etc/rc.local
#!/usr/bin/env bash
chmod 777 /sys/class/backlight/intel_backlight/brightness

