# config/systemd/user/backup.timer
[Unit]
Description=Backup timer

[Timer]
OnCalendar=daily
Persistent=false

[Install]
WantedBy=default.target
