# config/systemd/user/offlineimap-oneshot@.timer
[Unit]
Description=Offlineimap Query Timer for account %i

[Timer]
OnBootSec=1m
OnUnitInactiveSec=30m

[Install]
WantedBy=default.target
