# config/systemd/user/vdirsyncer.timer
[Unit]
Description=Vdirsyncer timer

[Timer]
OnBootSec=1m
OnUnitInactiveSec=30m

[Install]
WantedBy=default.target
