# etc/udev/rules.d/99-desktop.rules
ACTION=="change", KERNEL=="card0", SUBSYSTEM=="drm", RUN+="/bin/systemctl start kmarc-monitor"
